# supervisor.ini
[supervisord]
nodaemon=true

[program:example]
environment=
  SOME_VAR="%(ENV_SOME_VAR)s"
command=python /Users/mauro/Desktop/run.py 
