# supervisor.ini
[supervisord]
nodaemon=true

[program:run]
environment=
  SOME_VAR="%(ENV_SOME_VAR)s"
command=python /tmp/run.py 
